# systemd/ufo-galaxy-health.service
[Unit]
Description=UFO Galaxy Health Monitor
After=ufo-galaxy.service
Wants=ufo-galaxy.service

[Service]
Type=simple
User=ufo-galaxy
Group=ufo-galaxy
WorkingDirectory=/opt/ufo-galaxy

ExecStart=/opt/ufo-galaxy/venv/bin/python -m health_monitor --watchdog --interval 30

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ufo-galaxy-health

[Install]
WantedBy=multi-user.target
